# etc/systemd/system/generate-catalog.service
[Unit]
Description=Generate INPX catalog CSV parts

[Service]
Type=oneshot
User=library
Group=library
WorkingDirectory=/home/library

# Run the generator with the venv Python
ExecStart=/home/library/venv/bin/python /home/library/generate_catalog.py

# Optional, but nice for logs
StandardOutput=journal
StandardError=journal
